# deploy/launchd/com.glitchcube.boot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  GlitchCube boot supervisor (user LaunchAgent).
  Install with bin/install-boot. Concrete bigbuddy paths — bigbuddy is the only prod host.
  Runs in the logged-in GUI session (auto-login required) so utmctl can see the VM.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.glitchcube.boot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/estiens/code/glitchcube-main/glitchcube_rails/bin/glitchcube-boot</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <!-- Crash recovery: KeepAlive=true so launchd restarts the whole stack (glitchcube-boot
       -> foreman -> Puma + SolidQueue) whenever it exits. A post-boot Rails crash self-heals
       instead of leaving the cube dark until the next power-cycle.
       Because a plain TERM would just be resurrected, `bin/glitchcube-ctl stop` takes it down
       with `launchctl bootout` (job removed -> stays down for maintenance; the plist is still
       on disk so a real reboot starts it again via RunAtLoad). `restart` sends a graceful TERM
       and lets KeepAlive respawn it. glitchcube-boot reaps any orphaned Puma/SolidQueue from a
       hard-killed previous instance before re-exec'ing foreman, so restarts never double-run
       workers or lose port 4567. -->
  <key>KeepAlive</key>
  <true/>

  <!-- Minimum seconds between (re)launches. Softens the log/CPU cost of a permanent-error
       crash loop (e.g. a bad migration): ~2 tries/min instead of launchd's default 6. Note
       boot.{out,err}.log are NOT rotated — if you ever see a persistent loop, fix the cause
       (it's failing loudly) or add newsyslog rotation for those two files. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>WorkingDirectory</key>
  <string>/Users/estiens/code/glitchcube-main/glitchcube_rails</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/estiens/.rbenv/shims:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/estiens</string>
    <!-- launchd has no locale, so Ruby defaults to US-ASCII and chokes on any
         UTF-8 it reads (this bit foreman parsing .env). Force a UTF-8 locale. -->
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>LC_ALL</key>
    <string>en_US.UTF-8</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/estiens/code/glitchcube-main/glitchcube_rails/log/boot.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/estiens/code/glitchcube-main/glitchcube_rails/log/boot.err.log</string>
</dict>
</plist>
